# firefox-rss.service
[Unit]
Description=Generate Firefox RSS

[Service]
Type=oneshot
ExecStart=/home/pbui/src/personal/firefox-rss/firefox-rss.run
